# common/.config/systemd/user/wallpaper-rotate.service
[Unit]
Description=Rotate wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/pick-wallpaper apply
