# macos/Library/LaunchAgents/homebrew.mxcl.mpd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>KeepAlive</key>
  <true/>
  <key>Label</key>
  <string>homebrew.mxcl.mpd</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/mpd</string>
    <string>--no-daemon</string>
    <string>/Users/nbr/.config/mpd/mpd.conf</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/opt/homebrew</string>
  <key>StandardErrorPath</key>
  <string>/opt/homebrew/var/log/mpd.log</string>
  <key>StandardOutPath</key>
  <string>/opt/homebrew/var/log/mpd.log</string>
</dict>
</plist>