# arch/.config/systemd/user/zoxide-seed-sync.service
[Unit]
Description=Regenerate zoxide seed.txt and seed the zoxide database

[Service]
Type=oneshot
ExecStart=%h/.local/bin/zoxide-seed-sync
